# TestApplication/app/src/main/assets/libreswan-setup/if-pre-up.d/test-sever-ubuntu.pem
-----BEGIN RSA PRIVATE KEY-----
MIIEpAIBAAKCAQEA0GMrJQ5iUQLv94FZpVDViJfLrTA3/hDurRljacvPz7dgMnCRxUkErVGVeyk6
PbkkKmRXMiuxMozVHbwqD4kKp0LRUfaVNZY2+yhxf+CigUjkpp6DrV7/hOmphzvQaIsh8dA9ZcZo
cVDGhhaTGM4UOS+oJk2X67lBTFUtk5mxiCUIRWAn2sb+5Uj8nu6ysy+kttLn2r2DcQL5ff0XLi5S
EpGS1PqehatdclwnkHdLoBjgDwHsPunfvlkQlffjJcnRHgyOoktC/934twJGj7vltyu7SEhLmFeQ
RwENc4Jm2wARbSyFH5lXF14PIM1xHCB9S/1ErwbHCs8PG4duY87eZQIDAQABAoIBAQCIvhpotZik
alYHRvUJexK0EE/R9Ms5R+km43sNyHU0a7OPRbkSCYktJlT9vZ7/M8dcnf1t6afccPfEDF5ZVWcN
PNZYf4OH52tqwi5AaKqaRbfs8iWcFbN8h0v8dq6ae3op0KkJaPnOKeAryxJ5e3lf0Tr473/aafVE
XryKl5ccq2kRHvnD0IKEKltr1m1psm2JM+A8TaStd2ScjRFDicdXUohx2fwjNLbd9lEjpoSzZasK
/cRvMh4VEJIbtARmatB9XPOfMfGOjdukAxRf2Dh4Q40fmqEIAaAJGm9bHjnNy2GXNzdxSwmQkKUx
YukubgM5RQtof0YAL9Y8Ui/xI/QhAoGBAO3YDHgd/Fas1tDB2XtWi+9ije6b1wnaT2sTlxjtflJI
6ZVD/hxg592wmEvwpeORkXwf6MWnr2GtD6kG8auBl1gqM3ND9AGbg3juqtfvahWj93gh6Svgswzi
VHAcb87Cg4vYJPrXRep/x2H8q5EtunQuT/ULRa0lQJnM+CXujWWNAoGBAOBLeqLkgEfpgbdwCltb
kaVp4a8vlyXqX5wJM/weM+PZw+YJ4pWBxWK6Ht6eZbZcPxK74asL74uLJ4m87ncPVZUmXeht+DJb
tSQm50LkObzhXHuoqvh8Ea/yEv539t3QjsPMFogF/CzGSv/32U6IPEvyP++zgFf7TAvbrJZLeso5
AoGAU4sL9jUyKY/UdMdiiPxBcvh/OkzrH77bHpVnAY6KTITn51JWs8kmKKybEjGZU3ISVDMxJBxe
ckTDt2+E++UtL/7Cy74Wu2nQqG3M/MRP0wIhLsVajdsDiYcGecrDSzDv2/M0Khfrqw7vDfAeo4QY
di6+TyiCdOsc4ekpUMmWSRECgYA65Ja23riwKSfeOpT9l/WOMVcI4L4XPosMNX5o3KVTVoKRjFbq
WsCxQ52HZp/eWipViO+gMeocAdXH3Mvghy45k3Fk7jsVPe86JtziO3rNlX7rGUaLKEicojtKhcVU
SMazsorX/HpSZ16aB0IP1MXAS0NudPtmHp3BH1UZwYfkQQKBgQCXzwP9OHPADZHrWwHcw3TwEzKb
AOgDClrntAHK1rGJL/PXiKUsGtSGuCh/5XmHIAV0b5FVhHnn/pPPZ7VmVEwkPfy/8B0clKQ9jYg1
pSHEH4zO6rZtxyvRrYwkpueH8cC8Y6KNX5dp7UX2LGII5BTvFsXh4SSgb7ltnG6k19GJ8w==
-----END RSA PRIVATE KEY-----